# letsencrypt-created-certificate@.path
[Unit]
Description= Watches updates to LetsEncrypt certificates for: %i
Documentation= man:systemd.path

[Path]
PathChanged=/etc/letsencrypt/archive/%i/

[Install]
WantedBy=multi-user.target
