# media-enabled-vhosts@.path
[Unit]
Description= Watches changes to Media enabled VirtualHost: %i.
Documentation= man:systemd.path

[Path]
PathChanged=/var/www/vhosts/system/%i/conf/nginx.conf

[Install]
WantedBy=multi-user.target
