# php-5.6-etc@.path
[Unit]
Description= Watches changes to PHP FPM Pool definition: %p/%i.
Documentation= man:systemd.path

[Path]
PathChanged=/opt/plesk/%P/php-fpm.d/%i.conf

[Install]
WantedBy=multi-user.target
