# wp_db_export@examplepath.service.d/custompath.conf
[Service]
Environment=PUBLIC_PATH=wordpress/public/cms
