# wp_rsync@.service.d/10-override-example.conf
[Service]
Environment=REMOTE_HOST=user@0.tcp.ngrok.io
Environment="SSH_SPEC=ssh -p 22022 -o StrictHostKeyChecking=no -o UserKnownHostsFile=/dev/null"
